# files/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin tom --noclear %I $TERM
ExecStartPost=/usr/bin/su - tom -c '/usr/bin/uwsm start'
